# services/mcp-router.service
[Unit]
Description=MCP Email Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=%i
WorkingDirectory=/home/%i/mcp-infrastructure

# Load environment from .env
EnvironmentFile=/home/%i/mcp-infrastructure/.env

# Run with uv
ExecStart=/home/%i/.local/bin/uv run mcp-email-server

Restart=on-failure
RestartSec=5

StandardOutput=journal
StandardError=journal
SyslogIdentifier=mcp-router

[Install]
WantedBy=multi-user.target
